# crontab
@reboot cd /home/ezequiel/.vscode-oss && git pull origin master
*/15 * * * * cd /home/ezequiel/.vscode-oss && git add . && git commit -m "update-date +\%F--\%H:\%M" && git push origin master

